# rhel/dicomweb-proxy.service
[Unit]
Description=DICOM DIMSE to DICOMweb Proxy Server
Documentation=https://github.com/your-org/js-dicomweb-proxy
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=dicomweb
Group=dicomweb
WorkingDirectory=/opt/dicomweb-proxy
ExecStart=/opt/dicomweb-proxy/dicomweb-proxy-linux /opt/dicomweb-proxy/config/config.jsonc
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dicomweb-proxy

# Signal handling and shutdown configuration
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
SendSIGKILL=yes

# Security settings - relaxed for reliability
NoNewPrivileges=false
PrivateTmp=false
PrivateDevices=false
ProtectHome=false
ProtectSystem=false
# Removed ReadWritePaths and ReadOnlyPaths restrictions

# Allow binding to privileged ports without running as root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target